# deploy/backup/metagraphed-pg-backup.timer
# See metagraphed-pg-backup.service for setup. Daily cadence matches the RPO
# rationale in backup-postgres.sh (chain data is re-derivable, so ~daily is
# plenty — no PITR needed). RandomizedDelaySec avoids a fixed-second stampede
# against other box cron (apt/mdcheck/etc.); Persistent=true catches up a
# run missed by a reboot instead of silently skipping it.
[Unit]
Description=Daily metagraphed indexer Postgres -> R2 backup

[Timer]
OnCalendar=*-*-* 04:00:00 UTC
RandomizedDelaySec=15min
Persistent=true

[Install]
WantedBy=timers.target
