# deploy/backup/metagraphed-pg-backup.service
# Bare-metal equivalent of deploy/backup.railway.json — same image (this
# directory's Dockerfile + backup-postgres.sh, unmodified), triggered by the
# companion .timer instead of Railway's managed cron. For a self-hosted
# indexer box, not Railway.
#
# One-time setup on the box:
#   1. Copy this directory's Dockerfile + backup-postgres.sh to the box and
#      `docker build -f deploy/backup/Dockerfile -t metagraphed-pg-backup .`
#      (build context = the copied deploy/backup/'s parent, matching the
#      Dockerfile's `COPY deploy/backup/backup-postgres.sh` path).
#   2. Write the env vars this script needs (DATABASE_URL pointed at the
#      box's local Postgres, R2_BUCKET, R2_ENDPOINT, AWS_ACCESS_KEY_ID,
#      AWS_SECRET_ACCESS_KEY, BACKUP_PREFIX) to a root-only file, e.g.
#      /opt/metagraphed-indexer/backup/backup.env (chmod 600) — NEVER commit
#      that file; update ExecStart's --env-file path below to match.
#   3. Copy both units to /etc/systemd/system/, then:
#      systemctl daemon-reload && systemctl enable --now metagraphed-pg-backup.timer
#   4. Verify once with `systemctl start metagraphed-pg-backup.service` and
#      confirm the object lands in R2 before trusting the timer alone.
[Unit]
Description=metagraphed indexer Postgres -> R2 backup (mirrors deploy/backup.railway.json)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStartPre=-/usr/bin/docker rm -f metagraphed-pg-backup
ExecStart=/usr/bin/docker run --rm --name metagraphed-pg-backup --network host --env-file /opt/metagraphed-indexer/backup/backup.env metagraphed-pg-backup:latest
